# robot_bringup/systemd/robot_vision.service
[Unit]
Description=RoboMaster Robot System - Vision (Auto Aimer)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/nvidia/NJU-algorithm/vision_task/armor_task/build
Environment="ROS_DOMAIN_ID=42"
Environment="RMW_IMPLEMENTATION=rmw_cyclonedds_cpp"
Environment="LD_LIBRARY_PATH=/usr/local/opencv-4.12/lib:/usr/local/lib:/opt/ros/humble/lib:/usr/local/cuda/lib64"
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && cd /home/nvidia/NJU-algorithm/vision_task/armor_task/build && ./auto_aimer_test'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=robot_system.target
